# distros/debian/shutorial.cron.d
*/5 * * * * /usr/sbin/shutorial-admin schroot-users >/dev/null 2>&1
